# gopald.service
[Unit]
Description=Music Listening Tracker Daemon
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gopald --foreground
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=RUST_LOG=info
Environment=XDG_RUNTIME_DIR=%t

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share
ReadWritePaths=%h/.config

# Resource limits
MemoryMax=100M
CPUQuota=10%

[Install]
WantedBy=default.target
